# Strong Seaweed/blue-weave-aqua-main/server/.systemd/blueweave-backend.service
[Unit]
Description=Blue Weave Node.js Backend (Express + MongoDB)
After=network.target mongodb.service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/start-blueweave-backend.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
