# Strong Seaweed/blue-weave-aqua-main/server/.systemd/blueweave-backend.service
[Unit]
Description=Blue Weave Node.js Backend (Express + MongoDB)
After=network.target mongodb.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server
Environment="NODE_ENV=production"
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/_backend_stdout.log
StandardError=append:/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/_backend_stderr.log

[Install]
WantedBy=multi-user.target
