# Strong Seaweed/blue-weave-aqua-main/server/.systemd/python-agents.service
[Unit]
Description=Python Agent Gateway (FastAPI + LangChain/CrewAI)
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/start-python-agents.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
