# Strong Seaweed/blue-weave-aqua-main/server/.systemd/python-agents.service
[Unit]
Description=Python Agent Gateway (FastAPI + LangChain/CrewAI)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory="/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/agents_python"
Environment="PATH=/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/agents_python/.venv/bin:/usr/bin"
ExecStart="/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/agents_python/.venv/bin/uvicorn" agent_gateway:app --host 0.0.0.0 --port 8101
Restart=always
RestartSec=10
StandardOutput=append:"/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/agents_python/_agent_gateway_stdout.log"
StandardError=append:"/home/ubuntu/Strong Seaweed/Strong Seaweed/blue-weave-aqua-main/server/agents_python/_agent_gateway_stderr.log"

[Install]
WantedBy=multi-user.target
