# Strong Seaweed/Strong-Seaweed-main/.systemd/ml-api.service
[Unit]
Description=ML Seaweed Prediction APIs (Kappaphycus & Species)
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/start-ml-api.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
